# resources/monitord.service
[Unit]
Description=System monitoring daemon for COSMIC Observatory

[Service]
ExecStart=/usr/bin/monitord

[Install]
WantedBy=multi-user.target
